# udev/99-homey-slave.rules
# Homey Pro USB detection rules
# When a Homey Pro is connected via USB, it creates a USB Ethernet interface.
# This rule triggers the broker to register the device.
#
# ACTION: udev calls our registration script on USB add/remove events.
# ATTRS: Update vendor/product IDs when actual Homey Pro is connected (use `lsusb`).

# Homey Pro USB Ethernet adapter (placeholder IDs — update with real values)
ACTION=="add", SUBSYSTEM=="net", ATTRS{idVendor}=="0bda", ATTRS{idProduct}=="8152", \
    RUN+="/opt/pve-net-broker/scripts/on-usb-event.sh add %k"

ACTION=="remove", SUBSYSTEM=="net", ATTRS{idVendor}=="0bda", ATTRS{idProduct}=="8152", \
    RUN+="/opt/pve-net-broker/scripts/on-usb-event.sh remove %k"
